# saltstack/salt/helloworld-gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=0.0.0.0:8113

[Install]
WantedBy=sockets.target
