# saltstack/salt/gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ pillar['app_dir'] }}/gunicorn.socket

[Install]
WantedBy=sockets.target
